# aggregators/kagi-news/crontab
# Run Kagi News aggregator daily at 1 PM UTC (after Kagi updates around noon)
# Source environment variables exported by docker-entrypoint.sh
0 13 * * * . /etc/environment; cd /app && /usr/local/bin/python -m src.main >> /var/log/cron.log 2>&1

# Blank line required at end of crontab

