# aggregators/reddit-highlights/crontab
# Run Reddit Highlights aggregator every hour
# The Python script adds 0-10 min random jitter for anti-detection
# Effective polling interval: 60-70 minutes (varies each run)
# Source environment variables exported by docker-entrypoint.sh
0 * * * * . /etc/environment; cd /app && /usr/local/bin/python -m src.main >> /var/log/cron.log 2>&1

# Blank line required at end of crontab

